# roles/actions_runner/templates/actions-runner.service.j2
# Managed by Ansible - do not edit manually
[Unit]
Description=GitHub Actions Runner (kylar514/{{ runner_repo.name }}.{{ inventory_hostname }})
After=network.target ansible-keychain.service
Wants=ansible-keychain.service

[Service]
Type=simple
ExecStart={{ actions_runner_base_dir }}-{{ runner_repo.name }}/run.sh
WorkingDirectory={{ actions_runner_base_dir }}-{{ runner_repo.name }}
KillMode=process
KillSignal=SIGTERM
TimeoutStopSec=5min
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
